# systemd/cds-daemon.service
[Unit]
Description=CDS Zero-Trust Daemon
After=network.target

[Service]
ExecStart=/usr/local/bin/cds-daemon
Restart=always
User=root

# Redirect cosign/TUF cache to /tmp because of ProtectSystem=strict
Environment=COSIGN_CACHE=/tmp
Environment=TUF_ROOT=/tmp

# --- Security hardening ---
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

# Allow access to Docker socket for digest resolution and BoltDB storage
ReadWritePaths=/run/cds/ /var/run/docker.sock /var/lib/cds/
RestrictAddressFamilies=AF_UNIX AF_INET

# --- Correct way to create writable directories ---
RuntimeDirectory=cds
RuntimeDirectoryMode=0755
StateDirectory=cds
StateDirectoryMode=0750

[Install]
WantedBy=multi-user.target
