# systemd/cds-daemon.socket
[Unit]
Description=CDS Daemon Socket

[Socket]
ListenStream=/run/cds.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
